# unit_files/cleanup.timer
[Unit]
Description=Run Django Cleanup Command Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
